# pi/systemd/buzzer-daily.service
[Unit]
Description=Buzzer Pi Daily Health Summary

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/buzzer/pi
ExecStart=/usr/bin/python3 /home/pi/buzzer/pi/daily_summary.py
StandardOutput=journal
StandardError=journal
