# pi/systemd/buzzer-bot.service
[Unit]
Description=Buzzer Pi Telegram Bot
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/buzzer/pi
ExecStart=/usr/bin/python3 /home/pi/buzzer/pi/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
